# prometheus/prometheus.service
[Unit]
Description="Prometheus Monitoring"
Documentation=https://prometheus.io/
Requires=network-online.target
After=network-online.target
ConditionFileNotEmpty=/etc/prometheus/prometheus.yml

[Service]
User=prometheus
Group=prometheus
ExecStart=/usr/local/bin/prometheus --config.file /etc/prometheus/prometheus.yml --storage.tsdb.path=/var/prometheus
ExecReload=killall -HUP prometheus
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target